# sys-auth/validity-sensor/files/vcsFPService.service.in
## Copyright 2016 Jan Chren (rindeal)
## Distributed under the terms of the GNU General Public License v2

## vcsFPService.service
## The name vcsFPService was chosen based on vcsFPService.exe Windows service.

## This file was ported from:
##  - etc/init.d/vcsFPServiceDaemon

[Unit]
Description=Validity Fingerprint Service
# load before bringing up console so that FP authentication works
Before=multi-user.target

[Service]
# delete single instance lock file
ExecStartPre=/bin/rm -f /var/run/ValidityService.pid
# clean `/tmp/vcsSemKey_* junk`
ExecStartPre=/usr/bin/find /tmp -maxdepth 1 -type f -name "vcsSemKey_*" -delete
# clean `/tmp/CH_* junk`
ExecStartPre=/usr/bin/find /tmp -maxdepth 1 -type p -name "CH_*" -delete
# these files cannot be created without root permissions
ExecStartPre=/bin/touch /var/run/ValidityService.pid
ExecStartPre=/bin/chown @VFS_USER@:@VFS_GROUP@ /var/run/ValidityService.pid

User=@VFS_USER@
Group=@VFS_GROUP@
# libvfsFprintWrapper.so will still be creating files in /tmp with explicit 0666 perms
UMask=@VFS_UMASK@
# permissions applied to ExecStart only
PermissionsStartOnly=true

# daemon creates:
#  - `/tmp/vcsSemKey_*` files (pipes, IPC)
#  - /tmp/CH_* (pipes)
#  - /var/run/ValidityService.pid
#  - /etc/ValidityPersistentData
# daemon tries to access:
#  - /home/validity/Br_4_5_1xx/usdk/vcsFPService/../lib/linux/x64/Release/*
#  - /etc/localtime
#  - /usr/lib64/locale/locale-archive
#  - /etc/udev/udev.conf
#  - /run/udev/control
#  - /run/udev/data/*
#  - /sys/subsystem
#  - /dev/bus/usb/*
#  - / /sys /sys/bus/ /sys/bus/usb/ /sys/bus/usb/devices
#  - /sys/bus/usb/*
#  - /sys/class
#  - /sys/devices/*
# daemon tries to execute:
#  - /bin/rm
#  - /bin/sh
#  - /usr/bin/pidof
InaccessiblePaths=/home
ReadOnlyPaths=/etc/
# `vcsFPService` communicates with `libvfsFprintWrapper.so` through tmp files
#PrivateTmp=true
# protect /boot and /usr
ProtectSystem=true

# access("/run/udev/control", F_OK) = 0
# socket(AF_NETLINK, SOCK_RAW|SOCK_CLOEXEC|SOCK_NONBLOCK, NETLINK_KOBJECT_UEVENT) = 5
RestrictAddressFamilies=AF_UNIX AF_NETLINK
# deny realtime scheduling
RestrictRealtime=true
NoNewPrivileges=true
SystemCallFilter=~ @clock @cpu-emulation @debug @keyring @module @mount @obsolete @privileged @raw-io

Type=forking
ExecStart=@VCSFPSERVICE_INST_DIR@/vcsFPService

# pidfile is empty, serves only as a single instance lock
#PIDFile=/var/run/ValidityService.pid
TimeoutSec=10s
Restart=on-abort

ExecReload=/bin/kill -HUP $MAINPID


[Install]
WantedBy=multi-user.target
