# net-p2p/rtorrent/files/rtorrentd.init.in
#!/sbin/openrc-run
# Author: Jan Chren <dev.rindeal+gentoo-rtorrent AT gmail.com>
# Released under the 3-clause BSD license.

RFL_PLUGINS="ssd screend"
.  @@EROOT@@usr/share/openrc-rfl/main || return 1

name='rTorrent daemon'
description='Daemonized rTorrent instance'

rfl_screend_command='/usr/bin/rtorrent'
rfl_screend_command_args="${rtorrent_args}"
rfl_screend_procname='rtorrent main'
rfl_screend_pidfile='/run/${RC_SVCNAME}.pid'

: ${retry:="TERM/5/TERM/5/TERM/2/KILL"}

rfl_ssd_args

depend() {
    need net
}

start_pre() {
    rfl_require_var user || return 1
}

start_post() {
    rfl_screend_start_post || return 1
}

stop_pre() {
    rfl_require_var user
    { read -r pid < "${rfl_screend_pidfile}"; } >/dev/null 2>&1
}

stop_post() {
    ebegin "Stopping remaining screen sessions"

    if  [ -n "$pid" ] && [ -d "/proc/${pid}" ] ;then
        eend 1 "rTorrent process with PID `${pid}` still exists"
        return 1
    elif rfl_screend_list_sessions | grep -q "^${rfl_screend_session_name}\$" ;then
        eend 1 "Sessions with name '${rfl_screend_session_name}' still exist for user '${user}'"
        return 1
    fi

    eend 0 "All right, everything seems to be dead"

    rm -f "$pidfile" >/dev/null 2>&1 || true
}
