# net-p2p/qbittorrent/files/qbittorrent@.service.in
[Unit]
Description=`%i`'s qBittorrent Daemon Service
Documentation=man:qbittorrent-nox
Documentation=https://github.com/qbittorrent/qBittorrent/wiki

After=network.target


[Service]
Type=simple

User=%i
WorkingDirectory=~

# daemon must not run as root
ExecStartPre=/usr/bin/test ${USER} != 'root'

ExecStart=@EXPAND_BINDIR@/qbittorrent-nox $QB_ARGS
TimeoutStartSec=30s

TimeoutStopSec=3m

Restart=on-abnormal

PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectControlGroups=true
NoNewPrivileges=true
ProtectKernelModules=true
MemoryDenyWriteExecute=true
RestrictRealtime=true


[Install]
WantedBy=multi-user.target
